# backend/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:frankenphp]
command=frankenphp run --config /etc/frankenphp/Caddyfile
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
stopasgroup=true
killasgroup=true

; In-container scheduler (supercronic, crontab at /etc/cron.d/ukolio): dispatches
; due Scheduled scripts every minute (scripts:tick) and runs the hourly due-date
; reminder (notifications:due-tick) and recurring-task (recurring-tasks:tick) ticks.
[program:cron]
command=supercronic /etc/cron.d/ukolio
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
stopasgroup=true
killasgroup=true

[program:amqp-consumer]
command=php /app/src/amqp-consumer.php
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
stopasgroup=true
killasgroup=true

; Dedicated consumer for the script-run queue. Runs the V8 sandbox, so it loads
; ext-v8js (-d extension=v8js.so) — scoped to this process so libnode/V8 stays
; out of FrankenPHP and the main amqp-consumer. The main consumer deliberately
; skips the ScriptRun queue (see amqp-consumer.php).
[program:script-worker]
command=php -d extension=v8js.so /app/src/script-worker.php
; Drop privileges: the V8 sandbox executes user-authored JS, so a V8 escape lands
; as this unprivileged UID (created in the Dockerfile), not root.
user=ukolio-script
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
stopasgroup=true
killasgroup=true
