# backend/docker/cron.d/kytario
SHELL=/bin/sh
0 * * * *  php /app/bin/console notifications:due-tick
